# conf/supervisor-hbase.conf
[program:hbase]
command=/data/start-hbase.sh
stopwaitsecs=60
priority=10
